# acm-django/deploy/supervisor/acm-django.conf
[program:acm-django]
command = /usr/local/bin/acm-django-start
user = acm
stdout_logfile = /var/log/acm-django/stdout.log
stderr_logfile = /var/log/acm-django/stderr.log

